# root/etc/s6-overlay/s6-rc.d/svc-zim-updater/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid abc /update.sh &
wait

# Shut the pod down, we're done
echo "We're done here. Shutting down..."
poweroff
